# config/rc.local
#!/bin/bash -e

/usr/bin/mongod -f /etc/mongod.conf
